# lvm/aliyun/rc.local
#!/bin/sh

/etc/squids/create_lvm_in_aliyun.sh