# lvm/azure/rc.local
#!/bin/sh

/etc/squids/create_lvm_in_azure.sh